# files/mariadb.ini
[program:mariadb]
command=/usr/bin/pidproxy /var/run/mariadb/mariadb.pid /usr/bin/mysqld_safe --datadir='/var/lib/mysql'
autostart = true
